# .config/power-mode/99-auto-cpufreq-restart.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/usr/bin/systemctl --no-block restart auto-cpufreq.service"
